# services/learning-loop/deploy/gap-analyzer.service
# systemd service for gap analysis — MDJ_SERVER
# System Synergy — ZRD-SYN-008
#
# Deploy to: /etc/systemd/system/gap-analyzer.service
# Then: sudo systemctl daemon-reload && sudo systemctl enable gap-analyzer.timer

[Unit]
Description=RPI Gap Analyzer — Daily knowledge gap analysis
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=jdm
WorkingDirectory=/home/jdm/Projects/toMachina/services/learning-loop
Environment="NODE_ENV=production"
Environment="GOOGLE_APPLICATION_CREDENTIALS=/home/jdm/Projects/dojo-warriors/mdj-agent/sa-key.json"
EnvironmentFile=/home/jdm/Projects/dojo-warriors/mdj-agent/.env
ExecStart=/usr/bin/node dist/gap-analyzer.js
StandardOutput=journal
StandardError=journal
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
