# services/learning-loop/deploy/session-extract.timer
# systemd timer for session extraction — MDJ_SERVER
# Sprint 003 Learning Loop 2.0 — TRK-14168
#
# Deploy to: /etc/systemd/system/session-extract.timer
# Then: sudo systemctl daemon-reload && sudo systemctl enable --now session-extract.timer

[Unit]
Description=Weekly session extract to GCS (Saturday 11pm)

[Timer]
OnCalendar=Sat *-*-* 23:00:00
Persistent=true

[Install]
WantedBy=timers.target
