# services/learning-loop/deploy/wire-brain-sync.service
# systemd service for warrior brain sync — MDJ_SERVER
# System Synergy — ZRD-SYN-010
#
# Deploy to: /etc/systemd/system/wire-brain-sync.service
# Then: sudo systemctl daemon-reload && sudo systemctl enable wire-brain-sync.timer

[Unit]
Description=RPI Wire Brain Sync — Daily cross-warrior brain synchronization
After=network-online.target entity-extractor.service
Wants=network-online.target

[Service]
Type=oneshot
User=jdm
WorkingDirectory=/home/jdm/Projects/toMachina/services/learning-loop
Environment="NODE_ENV=production"
Environment="GOOGLE_APPLICATION_CREDENTIALS=/home/jdm/mdj-agent/sa-key.json"
EnvironmentFile=/home/jdm/Projects/dojo-warriors/mdj-agent/.env
ExecStart=/usr/bin/node dist/brain-export.js
StandardOutput=journal
StandardError=journal
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
