# services/learning-loop/claude-md-tracker.service
[Unit]
Description=CLAUDE.md Evolution Tracker — Weekly diff summary to knowledge_entries
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=jdm
WorkingDirectory=/home/jdm/mdj-agent
Environment="NODE_ENV=production"
Environment="GOOGLE_APPLICATION_CREDENTIALS=/home/jdm/mdj-agent/sa-key.json"
Environment="DOJO_WARRIORS_PATH=/home/jdm/Projects/dojo-warriors"
ExecStart=/usr/bin/node /home/jdm/mdj-agent/services/learning-loop/dist/claude-md-tracker.js
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
